# mpd-cdda-playlist.rules
KERNEL=="sr[0-9]", ACTION=="change", ENV{DISK_MEDIA_CHANGE}=="1", RUN+="/usr/bin/mpd-cdda-playlist.sh /var/lib/mpd/playlists/CDDA.m3u"
